# ops/systemd/arctic-rss-backup.timer
[Unit]
Description=Daily Arctic RSS verified PostgreSQL backup

[Timer]
OnCalendar=*-*-* 03:20:00 UTC
Persistent=true
RandomizedDelaySec=15m
Unit=arctic-rss-backup.service

[Install]
WantedBy=timers.target
