# ops/systemd/arctic-rss-build-cache-prune.timer
[Unit]
Description=Prune unused Arctic RSS Docker build cache every 72 hours

[Timer]
OnUnitActiveSec=72h
Persistent=true
Unit=arctic-rss-build-cache-prune.service

[Install]
WantedBy=timers.target
