# ops/systemd/arctic-rss-build-cache-prune.timer
[Unit]
Description=Weekly Arctic RSS Docker build-cache cleanup

[Timer]
OnCalendar=Sun *-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=30m
Unit=arctic-rss-build-cache-prune.service

[Install]
WantedBy=timers.target
